# wilson/ai-assistant.service
[Unit]
Description=Wilson AI-Assistent (Telegram)
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/reinhard/.openclaw/scripts/ai_assistant.py
Restart=on-failure
RestartSec=15
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/home/reinhard/.openclaw/secrets.env
Environment=TELEGRAM_CHAT_ID=8620231031
Environment=DEEPSEEK_MODEL=deepseek-chat
Environment=ENZYME_URL=http://192.168.86.195:11180
Environment=POLL_TIMEOUT=20
Environment=DB_PATH=/home/reinhard/.openclaw/ai_assistant.db

[Install]
WantedBy=default.target
